# .docker/build/php/config/supervisor/arcanesoft-worker.conf
[program:arcanesoft-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3
numprocs=4
stdout_logfile= /var/log/arcanesoft-worker.log
stderr_logfile_maxbytes=1MB
autostart=true
autorestart=true
startretries=0
